# ubuntu/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; path to your socket file
chmod = 0777

[supervisord]
nodaemon=true

[program:nginx]
command = /usr/sbin/nginx -c /etc/nginx/nginx.conf
