# ubuntu/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command = service nginx start

[program:sshd]
command = service ssh start
